# conf/vault/vault-login.timer
[Unit]
ConditionFileNotEmpty=/var/lib/vault/yubikey-seal.json

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
